# gunicorn.conf
[program:gunicorn]
command=gunicorn -c gunicorn_config.py wsgi:application
directory=/home/vagrant
autostart=true
autorestart=true
redirect_stderr=true
user=root
logpath=/var/log/gunicorn.log