# gunicorn.conf
[program:gunicorn]
command=sudo /opt/env/bin/gunicorn -c gunicorn_config.py xbrowse_server.wsgi
directory=/vagrant/code/xbrowse-web
autostart=true
autorestart=true
redirect_stderr=true