# mongodb.conf
[program:mongodb]
command=mongod --dbpath /data/db
user=root
autostart=false