# mongodb.conf
[program:mongodb]
command=mongod
user=root
autostart=true